# ashell/dot-config/systemd/user/ashell.service
[Unit]
Description=aShell
# Ensures it starts after the GUI session is ready
PartOf=hyprland-session.target
After=pipewire.service pipewire-pulse.service

[Service]
ExecStart=%h/.local/bin/ashell
Type=simple
Restart=always

[Install]
WantedBy=hyprland-session.target

